# systemd/ppod.service
[Unit]
Description=PPOD Webhook Server
After=network.target

[Service]
Type=simple
User=ppod
WorkingDirectory=/opt/ppod

ExecStart=/usr/bin/node /opt/ppod/ppod.js

Restart=always
RestartSec=3

Environment=NODE_ENV=production
EnvironmentFile=/etc/ppod/ppod.conf

# Logging (optional)
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target